# media-server/crossline-maintenance.timer
[Unit]
Description=Run Crossline API maintenance every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=15s
Persistent=true
Unit=crossline-maintenance.service

[Install]
WantedBy=timers.target
